# setup/codewars-runner.conf
start on runlevel [12345]
respawn

start script
    sudo sh /codewars-runner/setup/update.sh
end script